# hwp_encoder/sobonelib-hwp-encoder.service
[Unit]
Description=sobonelib hwp_encoder process

[Service]
ExecStart=
User=
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
